# proxy/proxy.service
[Unit]
Description=Message Broker Proxy Server
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=messagebroker
Group=messagebroker
WorkingDirectory=/opt/message_broker_proxy/proxy
Environment="PATH=/opt/message_broker_proxy/venv/bin"
EnvironmentFile=-/opt/message_broker_proxy/.env
Environment="PYTHONPATH=/opt/message_broker_proxy"

# Start command
ExecStart=/opt/message_broker_proxy/venv/bin/uvicorn app:app \
    --host 0.0.0.0 \
    --port 8001 \
    --ssl-keyfile certs/proxy.key \
    --ssl-certfile certs/proxy.crt \
    --ssl-ca-certs certs/ca.crt \
    --workers 4 \
    --log-level info

# Restart policy
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/message_broker_proxy/proxy/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target

